# submission/week-01/07.timer
[Unit]
Description=Run Hello World Service Every Minute

[Timer]
Unit=07.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target